# launchd/com.example.podcast-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily fetch job. Edit REPLACE_ME_USERNAME and the install path before use.
  Install: cp to ~/Library/LaunchAgents/ and `launchctl load` it.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.podcast-fetch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/REPLACE_ME_USERNAME/podcast-digest-template/launchd/wrapper_fetch.sh</string>
    </array>

    <!-- Daily at 6:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME_USERNAME/podcast-digest-template</string>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME_USERNAME/podcast-digest-template/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME_USERNAME/podcast-digest-template/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
